# cron/1.hour.cron
#!/bin/bash
source ve/bin/activate && python newsconnector/manage.py newsfeed_cron && deactivate
